# webmenu-editor-1.0/etc/profile.d/webmenu-update-hook.sh~

/usr/share/bin/webmenu-update.sh
logger webmenu-update "$(groups)" ${HOME}
